# ssh-downloader.service.example
[Unit]
Description=SSDOV SSH Downloader
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/YOURUSER/ssh-downloader
Environment=SSHDOWN_PASSWORD=change-this-password
# Or use public-key auth instead:
# Environment=SSHDOWN_AUTHORIZED_KEYS=/home/YOURUSER/.ssh/authorized_keys
ExecStart=/home/YOURUSER/ssh-downloader/ssdov -addr :2222 -root /srv/downloads -user download
# To use one folder for both browsing/downloads and SCP uploads, omit -root and add:
# ExecStart=/home/YOURUSER/ssh-downloader/ssdov -addr :2222 -user download -upload /home/YOURUSER/Downloads
# To enable legacy SCP uploads to a different folder, add: -upload /home/YOURUSER/Downloads
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
